# crontab
* * * * * git -C /node-cron pull && cp /node-cron/index.html /etc/nginx/html/ && nginx -s reload >> /var/log/cron.log 2>&1
