# crontab
* * * * * echo "git pull https://github.com/ioforce/node-cron.git . && copy index.html /etc/nginx/html/ && nginx -s reload" >> /var/log/cron.log 2>&1
